# src/main/bin/http-demo.ini
[program:http-demo]
command=/opt/http-demo/startup.sh
directory=/opt/http-demo
user=root
autostart=true
autorestart=true
stopsinal=INT
stdout_logfile=/var/log/http-demo.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
